# ansible/main/roles/postgres/files/postgres-backups/pg_backup.service
[Unit]
Description=Backup postgres cluster
Wants=pg_backup.timer

[Service]
ExecStart=/home/player/postgres-backup.sh

[Install]
WantedBy=multi-user.target
